# scripts/php-fpm.ini
[program:php-fpm]
command=/usr/local/bin/php-fpm.sh start
user=root
daemonize = no
